# etc/cron.d/poll-changes
PATH=/usr/local/bin:/usr/bin:/bin
SHELL=/bin/sh

*/15 * * * * ubuntu /usr/local/bin/poller.sh
